# docker/service/node/run
#!/bin/bash
exec node /var/www/gitcolab/bin/server.js